# scripts/activity-bot/activity-bot.timer
[Unit]
Description=Run Activity Bot daily at 10:00 UTC

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
